# systemd/backup.service
[Unit]
Description=Backup cronjob

[Service]
Type=oneshot
ExecStart=/root/backup/run

[Install]
WantedBy=multi-user.target
